# playbooks/roles/config/lxc_docker_environment/templates/overmind-postgres-backup.service.j2
[Unit]
Description=Overmind Postgres verified logical backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/overmind-postgres-backup.env
ExecStart=/usr/local/sbin/overmind-postgres-backup
